# deploy/systemd/templates/thesis-app.service.tmpl
[Unit]
Description=Thesis Flask API and Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__PROJECT_ROOT__
EnvironmentFile=-__PROJECT_ROOT__/.env
EnvironmentFile=__PROJECT_ROOT__/deploy/systemd/thesis.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc 'exec "$PYTHON_BIN" pi/app.py'
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
